# jamf/Extension Attributes/AI Software Inventory/io.github.ai-software-inventory.collector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- SPDX-FileCopyrightText: 2026 Robert Flanagan and macadmin-toolbox contributors
     SPDX-License-Identifier: MIT
     Reference copy only: install.sh generates the deployed plist at install time. -->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>io.github.ai-software-inventory.collector</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>--no-rcs</string>
		<string>/usr/local/ai-software-inventory/ai-inventory-collector.zsh</string>
	</array>
	<key>StartInterval</key>
	<integer>14400</integer>
	<key>RunAtLoad</key>
	<true/>
	<key>ProcessType</key>
	<string>Background</string>
	<key>StandardOutPath</key>
	<string>/var/log/ai-software-inventory.collector.log</string>
	<key>StandardErrorPath</key>
	<string>/var/log/ai-software-inventory.collector.log</string>
</dict>
</plist>
